# systemd/claude-backup.timer
[Unit]
Description=Claude Command Center — Nightly backup timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
