# systemd/claude-memory-dreaming.service
[Unit]
Description=claudhaus — Nightly memory dreaming sweep

[Service]
Type=oneshot
WorkingDirectory=%h/JjayFiles/claude-command-center
EnvironmentFile=-%h/JjayFiles/claude-command-center/.env
ExecStart=%h/JjayFiles/claude-command-center/.venv/bin/python memory/dreaming.py

StandardOutput=append:%h/JjayFiles/claude-command-center/logs/dreaming.log
StandardError=append:%h/JjayFiles/claude-command-center/logs/dreaming.log
