# systemd/claude-dashboard.service
[Unit]
Description=claudhaus — Dashboard (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/JjayFiles/claude-command-center/dashboard-ui
EnvironmentFile=-%h/JjayFiles/claude-command-center/.env
ExecStartPre=/home/jjay/.nvm/versions/node/v24.13.1/bin/npm run build
ExecStart=/home/jjay/.nvm/versions/node/v24.13.1/bin/npm start
Restart=on-failure
RestartSec=10

StandardOutput=append:%h/JjayFiles/claude-command-center/logs/dashboard.log
StandardError=append:%h/JjayFiles/claude-command-center/logs/dashboard.log

[Install]
WantedBy=default.target
